# jessie:openvswitch-2.5.0/supervisord.conf
[supervisord]
nodaemon=true

[program:ovs-ctl]
command=/usr/share/openvswitch/scripts/ovs-ctl start
